# ansible/roles/firewall/files/firehol.conf
description	"firehol firewall support"

start on runlevel [2345]
stop on runlevel [!2345]

exec firehol start
